# deploy/pixelrag-api.service
[Unit]
Description=PixelRAG Search API
After=network.target

[Service]
Type=simple
User=yichuan
WorkingDirectory=/home/yichuan/visrag
ExecStart=/home/yichuan/visrag/.venv/bin/pixelrag-serve \
    --index-dir /home/yichuan/visrag-data/search_index \
    --tiles-dir /home/yichuan/visrag-data/tiles \
    --articles-json /home/yichuan/visrag-data/articles.json \
    --device cpu \
    --port 30001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
